# 99_Repo_Exports/tick_flow_full/orderflow_services/deploy/systemd/exec-health-freeze-reconnect-nightly.timer
[Unit]
Description=Nightly ExecHealth reconnect chaos smoke timer

[Timer]
OnCalendar=*-*-* 03:17:00
RandomizedDelaySec=20m
Persistent=true
Unit=exec-health-freeze-reconnect-nightly.service

[Install]
WantedBy=timers.target
