# 99_Repo_Exports/systemd/trade-ml-calibration-health.timer
[Unit]
Description=Run Trade ML calibration health monitor daily

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
RandomizedDelaySec=30
AccuracySec=1min

[Install]
WantedBy=timers.target
















