# 99_Repo_Exports/python-worker/tick_flow_full/orderflow_services/deploy/systemd/trade-meta-cov-rollout-controller.timer
[Unit]
Description=Run trade meta-cov rollout controller every 10 minutes (latency-contract gated)

[Timer]
OnCalendar=*:0/10
Persistent=true
Unit=trade-meta-cov-rollout-controller.service

[Install]
WantedBy=timers.target
