# 99_Repo_Exports/python-worker/orderflow_services/deploy/systemd/trade-conf-score-guardrails-autopromo-controller.timer
[Unit]
Description=Run trade conf-score guardrails autopromo every 15 minutes (latency-contract gated)

[Timer]
OnCalendar=*:0/15
Persistent=true
Unit=trade-conf-score-guardrails-autopromo-controller.service

[Install]
WantedBy=timers.target
