# 99_Repo_Exports/python-worker/orderflow_services/deploy/systemd/trade-ofc-contextual-nightly-bundle.timer
[Unit]
Description=Run trade OFC contextual nightly bundle every day at 03:11 UTC

[Timer]
OnCalendar=*-*-* 03:11:00
Persistent=true
Unit=trade-ofc-contextual-nightly-bundle.service

[Install]
WantedBy=timers.target
