# 99_Repo_Exports/python-worker/orderflow_services/deploy/systemd/trade-ofc-contextual-rollout-controller.timer
[Unit]
Description=Run trade OFC contextual rollout controller every 10 minutes

[Timer]
OnCalendar=*:0/10
Persistent=true
Unit=trade-ofc-contextual-rollout-controller.service

[Install]
WantedBy=timers.target
