# 99_Repo_Exports/python-worker/orderflow_services/deploy/systemd/trade-tca-nightly-report-exporter.service
[Unit]
Description=Trade TCA nightly report exporter
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%E{TRADE_REPO_ROOT}
EnvironmentFile=%E{TRADE_ENV_FILE}
ExecStart=/usr/bin/env bash %E{TRADE_REPO_ROOT}/orderflow_services/deploy/systemd/run_trade_tca_nightly_report_exporter_v1.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
