# 99_Repo_Exports/deploy/systemd/trade-risk-engine-canary-report.timer
[Unit]
Description=Periodic trade risk engine canary report

[Timer]
OnBootSec=5m
OnUnitActiveSec=20m
Unit=trade-risk-engine-canary-report.service

[Install]
WantedBy=timers.target
