# 99_Repo_Exports/deploy/systemd/trade-execution-canary-report.timer
[Unit]
Description=Periodic trade execution canary report

[Timer]
OnBootSec=4m
OnUnitActiveSec=15m
Unit=trade-execution-canary-report.service

[Install]
WantedBy=timers.target
