# 99_Repo_Exports/deploy/systemd/trade-risk-decision-summary-refresh.timer
[Unit]
Description=Periodic refresh of trade risk decision summary

[Timer]
OnBootSec=6m
OnUnitActiveSec=30m
Unit=trade-risk-decision-summary-refresh.service

[Install]
WantedBy=timers.target
