# 99_Repo_Exports/deploy/systemd/trade-risk-decision-summary-refresh.service
[Unit]
Description=Fetch trade risk decision aggregate
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/trade
EnvironmentFile=-/etc/default/trade-risk-engine
ExecStart=/usr/bin/python3 /opt/trade/scripts/refresh_risk_decision_summary.py
