# 99_Repo_Exports/deploy/systemd/trade-risk-mismatch-summary-refresh.service
[Unit]
Description=Refresh risk mismatch summary materialized view
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/trade
EnvironmentFile=-/etc/default/trade-risk-p48
ExecStart=/usr/bin/python3 /opt/trade/scripts/refresh_risk_mismatch_summary.py
