# 99_Repo_Exports/deploy/systemd/trade-risk-mismatch-summary-refresh.timer
[Unit]
Description=Periodic refresh of risk mismatch summary

[Timer]
OnBootSec=6m
OnUnitActiveSec=15m
Unit=trade-risk-mismatch-summary-refresh.service

[Install]
WantedBy=timers.target
