# 99_Repo_Exports/deploy/systemd/trade-execution-retention.timer
[Unit]
Description=Nightly trade execution retention purge

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
Unit=trade-execution-retention.service

[Install]
WantedBy=timers.target
