# 99_Repo_Exports/etc/systemd/system/trade-of-nightly-meta-self-heal.timer
[Unit]
Description=Run meta cells self-heal nightly

[Timer]
OnCalendar=*-*-* 07:10:00
Persistent=true

[Install]
WantedBy=timers.target

