# 99_Repo_Exports/python-worker/systemd/trade-execution-auto-scrubber.service
[Unit]
Description=Trade execution auto-trigger checkpoint scrubber
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/trade
EnvironmentFile=-/etc/default/trade-execution-p33-autonomy
ExecStart=/usr/bin/python3 /opt/trade/scripts/auto_trigger_checkpoint_scrubber.py
