# 99_Repo_Exports/python-worker/systemd/trade-execution-auto-scrubber.timer
[Unit]
Description=Periodic trade execution auto-trigger checkpoint scrubber

[Timer]
OnBootSec=6m
OnUnitActiveSec=10m
Unit=trade-execution-auto-scrubber.service

[Install]
WantedBy=timers.target
