# 99_Repo_Exports/python-worker/systemd/trade-of-nightly-meta-enforce-ramp.service
[Unit]
Description=Trade OF nightly meta ENFORCE ramp-or-freeze proposal (worst-case stratified)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/trade/of_reports.env
WorkingDirectory=/home/alex/front/trade/scanner_infra/python-worker
ExecStart=/usr/bin/python3 -m tools.nightly_meta_enforce_ramp_or_freeze_bundle
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

