# 99_Repo_Exports/python-worker/tools/trade-of-weekly-bench.timer
[Unit]
Description=Run Trade OF weekly latency bench

[Timer]
OnCalendar=Sun *-*-* 06:40:00
Persistent=true

[Install]
WantedBy=timers.target

