# 99_Repo_Exports/python-worker/tools/trade-of-nightly-meta-unfreeze.timer
[Unit]
Description=Run meta unfreeze proposal nightly

[Timer]
OnCalendar=*-*-* 06:55:00
Persistent=true

[Install]
WantedBy=timers.target

