# 99_Repo_Exports/python-worker/tools/trade-of-regress.timer
[Unit]
Description=Run Trade OF regression report nightly

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true

[Install]
WantedBy=timers.target
