# 99_Repo_Exports/python-worker/tools/trade-of-regress.service
[Unit]
Description=Trade OF regression report (baseline inputs -> replay -> diff -> Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/trade/of_reports.env
WorkingDirectory=/home/alex/front/trade/scanner_infra/python-worker
ExecStart=/usr/bin/python3 -m tools.cron_of_reports --mode regress
