# 99_Repo_Exports/python-worker/tools/trade-of-nightly-meta-stage2-opt.timer
[Unit]
Description=Run meta Stage2 optimize nightly

[Timer]
OnCalendar=*-*-* 07:25:00
Persistent=true

[Install]
WantedBy=timers.target

