# 99_Repo_Exports/python-worker/tools/trade-cfg-drift.timer
[Unit]
Description=Run Trade config drift monitor hourly

[Timer]
OnCalendar=*-*-* *:15:00
Persistent=true

[Install]
WantedBy=timers.target

