# 99_Repo_Exports/reference/binance_orders/trade-execution-replay-slo-refresh.timer
[Unit]
Description=Periodic refresh of replay/rehydrate SLO summary

[Timer]
OnBootSec=8m
OnUnitActiveSec=15m
Unit=trade-execution-replay-slo-refresh.service

[Install]
WantedBy=timers.target
