# 99_Repo_Exports/reference/binance_orders/trade-execution-auto-repair.service
[Unit]
Description=Trade execution automated repair / quarantine policy
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/trade
EnvironmentFile=/opt/trade/.env
ExecStart=/usr/bin/python3 /opt/trade/scripts/automated_execution_repair_policy.py
