# 99_Repo_Exports/ops/systemd/edge-stack-v1-train-bundle.timer
[Unit]
Description=Edge Stack v1 nightly training bundle timer

[Timer]
OnCalendar=*-*-* 03:10:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
