# 99_Repo_Exports/ops/systemd/edge-stack-v1-shadow-status-exporter.service
[Unit]
Description=Edge Stack v1 shadow status Prometheus exporter
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/app
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/trade-edge-stack
ExecStart=/usr/bin/python -m orderflow_services.edge_stack_shadow_status_exporter_v1
Restart=always
RestartSec=5
User=trade
Group=trade

[Install]
WantedBy=multi-user.target
