# 99_Repo_Exports/ops/systemd/trade-ml-rollout-guard.service
[Unit]
Description=Trade ML rollout guard (freeze/unfreeze proposals)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/trade/of_reports.env
WorkingDirectory=/home/alex/front/trade/scanner_infra/python-worker
ExecStart=/usr/bin/python3 -m tools.ml_rollout_guard --since-min 60
StandardOutput=journal
StandardError=journal

