# 99_Repo_Exports/ops/systemd/trade-ml-rollback-proposer.timer
[Unit]
Description=Run ML rollback proposer hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

