# 99_Repo_Exports/ops/systemd/trade-ml-walkforward.timer
[Unit]
Description=Run ML walk-forward report daily

[Timer]
OnCalendar=*-*-* 05:40:00
Persistent=true

[Install]
WantedBy=timers.target
