# 99_Repo_Exports/ops/systemd/trade-ml-threshold-proposer.service
[Unit]
Description=Trade ML threshold proposer v1 (p_min proposals with Approve/Reject)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/trade/of_reports.env
WorkingDirectory=/home/alex/front/trade/scanner_infra/python-worker
ExecStart=/usr/bin/python3 -m tools.ml_threshold_proposer_v1

