# 99_Repo_Exports/ops/systemd/trade-ml-threshold-proposer-v2.timer
[Unit]
Description=Run ML threshold proposer v2 daily

[Timer]
OnCalendar=*-*-* 06:55:00
Persistent=true

[Install]
WantedBy=timers.target

