# 99_Repo_Exports/ops/systemd/trade-ml-nightly-train-report.timer
[Unit]
Description=Run ML nightly train+report daily

[Timer]
OnCalendar=*-*-* 05:10:00
Persistent=true

[Install]
WantedBy=timers.target

