# 99_Repo_Exports/ops/systemd/trade-ml-rollback-proposer.service
[Unit]
Description=Trade ML rollback proposer (propose restore *_prev, Approve/Reject)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/trade/of_reports.env
WorkingDirectory=/home/alex/front/trade/scanner_infra/python-worker
ExecStart=/usr/bin/python3 -m tools.ml_rollback_proposer_v1

