# 99_Repo_Exports/ops/systemd/trade-ml-rollout-guard-v2.timer
[Unit]
Description=Run ML rollout guard v2 hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

