# 99_Repo_Exports/deploy/systemd/user/notify-receiver.service
[Unit]
Description=OBI Notification Receiver
After=network.target redis.service book-analytics.service

[Service]
Type=simple
WorkingDirectory=/home/alex/front/trade/scanner_infra/python-worker
Environment=PYTHONPATH=/home/alex/front/trade/scanner_infra/python-worker
Environment=REDIS_URL=redis://127.0.0.1:6379/0
Environment=NOTIFY_STREAM=notify:telegram
Environment=NOTIFY_RECEIVER_PORT=8089
ExecStart=/usr/bin/python3 -m services.notify_receiver
Restart=always
RestartSec=5

[Install]
WantedBy=default.target

