# 99_Repo_Exports/deploy/systemd/user/orders-router.service
[Unit]
Description=Orders Router (callbacks -> orders queue)
After=network-online.target
Wants=network-online.target

[Service]
# %h = home directory
WorkingDirectory=%h/front/trade/scanner_infra/python-worker
ExecStart=/usr/bin/python3 %h/front/trade/scanner_infra/python-worker/services/orders_router.py

Restart=always
RestartSec=2

# Redis connection
Environment=REDIS_URL=redis://scanner-redis:6379/0

# Streams and queues
Environment=CALLBACKS_STREAM=bot:callbacks
Environment=ORDERS_QUEUE=orders:queue
Environment=SNAP_PREFIX=signal:snap:

# Consumer group
Environment=ORDERS_ROUTER_GROUP=orders-router-group
Environment=ORDERS_ROUTER_CONSUMER=orders-router-1

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

